# com.liquidityhunter.dailyupdate.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.liquidityhunter.dailyupdate</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/thkmacstudio1/LiquidityHunter/scripts/daily_update.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>1</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/thkmacstudio1/LiquidityHunter</string>

    <key>StandardOutPath</key>
    <string>/Users/thkmacstudio1/LiquidityHunter/logs/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/thkmacstudio1/LiquidityHunter/logs/launchd_error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
